# raspberry-pi-deployment/encoder.service
[Unit]
Description=Rotary Encoder Input Handler
Documentation=https://github.com/shreyashgupta/days-to-thing-tracker
After=kiosk.service
Requires=kiosk.service

[Service]
Type=simple
User=shreyash
Group=shreyash

# Environment for xdotool
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/shreyash/.Xauthority

# Working directory
WorkingDirectory=/home/shreyash

# Wait for X to be ready
ExecStartPre=/bin/sleep 3

# Start the encoder script
ExecStart=/usr/bin/python3 /home/shreyash/encoder.py

# Restart policy
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
