# raspberry-pi-deployment/kiosk.service
[Unit]
Description=Days Tracker Kiosk Browser
Documentation=https://github.com/shreyashgupta/days-to-thing-tracker
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=shreyash
Group=shreyash

# Environment
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/shreyash/.Xauthority
Environment=KIOSK_URL=https://days-tracker-server-deployment.reverse-python.ts.net/

# Working directory
WorkingDirectory=/home/shreyash

# Start X and then the kiosk script
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/xinit /home/shreyash/kiosk.sh -- :0 -nocursor

# Restart policy
Restart=on-failure
RestartSec=10

# Resource limits (Pi Zero 2 W has limited RAM)
MemoryMax=400M

[Install]
WantedBy=multi-user.target
