# systemd/kiosk-tracker.service
[Unit]
Description=Days Tracker Kiosk Display
After=network.target

[Service]
Type=simple
User=shreyash
Group=shreyash
WorkingDirectory=/home/shreyash/days-to-thing-tracker
Environment=PYTHONUNBUFFERED=1
Environment=KIOSK_DATA_DIR=/home/shreyash/days-to-thing-tracker/data
ExecStart=/home/shreyash/days-to-thing-tracker/venv/bin/python /home/shreyash/days-to-thing-tracker/python/main.py
Restart=always
RestartSec=5

# GPIO access
SupplementaryGroups=gpio

[Install]
WantedBy=multi-user.target
